# ami/mc-supervisor.conf
[program:minecraft]
user=ubuntu
directory=/home/ubuntu/
command=/home/ubuntu/launch_minecraft.sh
