# ami/watchdog.conf
[program:watchdog]
user=ubuntu
directory=/home/ubuntu/watchdog/
command=/home/ubuntu/watchdog/bin/python /home/ubuntu/watchdog/watchdog.py
autorestart=unexpected
